# crontab
#min hr  		md mo wkday command
30	 5,11,17,23	 *  *  *     curl 'http://localhost:8080/update'
0	 1  		 1  *  *	 curl 'http://localhost:8080/newmonth'